# deploy/systemd/plum-agent@.service
[Unit]
Description=Plum Agent - Node %i
Documentation=https://github.com/manxisuo/plum
After=network.target plum-controller.service
Requires=plum-controller.service

[Service]
Type=simple
User=manxisuo
Group=manxisuo
WorkingDirectory=/home/manxisuo/Plum

# %i 是实例名称（如nodeA、nodeB）
ExecStart=/home/manxisuo/Plum/agent-go/plum-agent

# 环境变量
Environment="AGENT_NODE_ID=%i"
Environment="CONTROLLER_BASE=http://127.0.0.1:8080"
Environment="AGENT_DATA_DIR=/home/manxisuo/Plum/data/agents/%i"

# 重启策略
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plum-agent-%i

[Install]
WantedBy=multi-user.target

